# units/systemd-machine-id-setup.service.in
#  This file is part of systemd/uselessd.
#
#  systemd is free software; you can redistribute it and/or modify it
#  under the terms of the GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Initialize Machine ID
Documentation=man:systemd-machine-id-setup(1)
DefaultDependencies=no
Wants=local-fs.target
Conflicts=shutdown.target
After=local-fs.target
Before=sysinit.target shutdown.target

[Service]
Type=oneshot
ExecStart=@rootbindir@/systemd-machine-id-setup
